# fragments/programs/davmail/home-config/davmail/.config/autostart/davmail.desktop
[Desktop Entry]
Type=Application
Exec=/run/current-system/sw/bin/davmail
Icon=davmail
Terminal=false
Categories=Network;Email;Java;

Name=DavMail
Comment=Mail agent for EWS
