# machines/thrud/home-config/autostart/.config/autostart/davmail.desktop
[Desktop Entry]
Type=Application
Exec=davmail
Icon=davmail
Terminal=false
Categories=Network;Email;

Name=DavMail
Comment=Mail agent for EWS
